# pkg/deb-root/DEBIAN/postinst
#!/bin/bash
set -e

# Create rqlite system user and group if they don't exist
if ! getent group rqlite >/dev/null 2>&1; then
    addgroup --system rqlite
fi

if ! getent passwd rqlite >/dev/null 2>&1; then
    adduser --system --no-create-home --home /var/lib/rqlite --shell /bin/false --group rqlite
fi

# Create data directory and set permissions
mkdir -p /var/lib/rqlite
chown rqlite:rqlite /var/lib/rqlite
chmod 755 /var/lib/rqlite

# Reload systemd daemon to recognize new service
if [ -d /run/systemd/system ]; then
    systemctl daemon-reload >/dev/null 2>&1 || true
fi

# Print installation message
echo "rqlite has been installed successfully."
echo "To start the rqlite service, run: sudo systemctl enable --now rqlite"
echo "Configuration can be modified in /etc/default/rqlite"
echo "Data directory: /var/lib/rqlite"